# supervisord.conf
[program:laravel]
user=www-data
command=/usr/bin/php /app/artisan queue:work
autostart=true
autorestart=true
startretries=10
kestrel=true
priority=999